# deploy/agent-server.service
[Unit]
Description=Appx Pi Agent Server
Documentation=https://github.com/appx-org/agent-server
After=network.target
Before=appx.service

[Service]
User=opencode
Group=opencode
UMask=0007

EnvironmentFile=/etc/appx/appx.env
Environment=HOME=/home/opencode

# Route provider HTTPS traffic through appx's egress CONNECT proxy on 127.0.0.1:9080.
# This keeps Pi model calls behind the same allowlist and request logging as the
# legacy OpenCode backend. Node.js needs NODE_USE_ENV_PROXY=1 before fetch/core
# HTTP clients honor HTTPS_PROXY/NO_PROXY.
Environment=NODE_USE_ENV_PROXY=1
Environment=HTTPS_PROXY=http://127.0.0.1:9080
Environment=NO_PROXY=localhost,127.0.0.1

Environment=AGENT_SERVER_MODE=multi
Environment=PROJECT_DIR=__APPX_PROJECTS_DIR__
Environment=SESSIONS_DIR=/home/opencode/.pi/agent/appx-default-sessions
Environment=AGENT_DIR=/home/opencode/.pi/agent
Environment=AGENTS_FILE=.pi/AGENTS.md
Environment=AGENT_SERVER_HOST=127.0.0.1
Environment=AGENT_SERVER_PORT=4001

WorkingDirectory=__APPX_PROJECTS_DIR__
ExecStart=/usr/local/bin/agent-server

Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
