# deploy/appx.service
[Unit]
Description=Appx — Agentic Application Proxy
Documentation=https://github.com/neuromaxer/appx
After=network.target

[Service]
User=appx
Group=appx

# Allows binding to port 443 without running as root.
# Remove these two lines if running on a non-privileged port (8443+).
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Ensure directories created by appx (project subdirs) have group-write so the
# agent service user can write files via the shared projects group.
UMask=0007

# Server-specific config lives in /etc/appx/appx.env (created by bootstrap).
# The env file sets APPX_DATA, APPX_HOST, APPX_PORT, etc.
EnvironmentFile=/etc/appx/appx.env

ExecStart=/usr/local/bin/appx

Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
